# flask/supervisord.conf
[supervisord]
nodaemon=true

[program:uwsgi]
directory=/www/backtend
command=/usr/sbin/uwsgi --ini /etc/uwsgi/uwsgi.ini --die-on-term
stdout_logfile=/log/uwsgi.access.log
stdout_logfile_maxbytes=0
stderr_logfile=/log/uwsgi.error.log
stderr_logfile_maxbytes=0

[program:nginx]
command=/usr/sbin/nginx
